# docker/processes.ini
[program:nginx]
command=/usr/sbin/nginx
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
